# actualizador/installer/debian/postinst
#!/bin/sh
set -e
systemctl daemon-reload 2>/dev/null || true
systemctl enable --now factupos-actualizador.timer 2>/dev/null || true
# correr una vez ya (en segundo plano)
systemctl start factupos-actualizador.service 2>/dev/null || true
exit 0
