# biometria-linux/packaging/deb-root/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop factupos-fingerprint-servicio.service 2>/dev/null || true
systemctl disable factupos-fingerprint-servicio.service 2>/dev/null || true
